# packaging/deb/daemon/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create daos user if it doesn't exist
        if ! id -u daos >/dev/null 2>&1; then
            useradd --system --no-create-home --shell /usr/sbin/nologin daos
        fi

        # Create directories
        mkdir -p /opt/daos
        mkdir -p /etc/daos
        mkdir -p /var/log/daos

        # Set ownership
        chown -R daos:daos /opt/daos
        chown -R daos:daos /var/log/daos

        # Create default config if it doesn't exist
        if [ ! -f /etc/daos/daemon.yaml ]; then
            cat > /etc/daos/daemon.yaml <<EOF
server:
  host: "127.0.0.1"
  port: 8080

database:
  path: "/opt/daos/daos.db"
EOF
            chown daos:daos /etc/daos/daemon.yaml
        fi

        # Enable and start service
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload
            systemctl enable daos-daemon.service
            systemctl start daos-daemon.service || true
        fi
        ;;
    abort-upgrade|abort-remove|abort-configure)
        ;;
esac

exit 0
