# packaging/deb/daemon/DEBIAN/prerm
#!/bin/bash
set -e

case "$1" in
    remove|purge)
        # Stop service
        if [ -d /run/systemd/system ]; then
            systemctl stop daos-daemon.service || true
            systemctl disable daos-daemon.service || true
        fi
        ;;
    upgrade|failed-upgrade)
        ;;
    remove)
        ;;
    purge)
        # Remove data directory
        rm -rf /opt/daos
        rm -rf /var/log/daos

        # Remove user
        userdel daos || true
        ;;
esac

exit 0
